# scripts/lunes_validator.service
#!/usr/bin/env bash
[Unit]
Description=Vaidador Lunes node
After=network.target

[Service]
Type=simple
Restart=always
RestartSec=1
User=root
ExecStart=/root/lunes-nightly/target/release/lunes-node \
                        --chain  /root/lunes-nightly/lunes-staging-raw.json -d /data/validator \
                        --name validator2 --validator --port 30333 --ws-port 9944 --rpc-port 9933
[Install]
WantedBy=multi-user.target
